# nyla.service
[Unit]
Description=Nyla — RAG Voice Assistant (Robotwala)
# Start only after network and sound subsystem are ready
After=network-online.target sound.target pulseaudio.service
Wants=network-online.target

[Service]
Type=simple

# ── Paths ─────────────────────────────────────────────────
WorkingDirectory=/home/acrossd/Desktop/nyla
ExecStart=/home/acrossd/Desktop/nyla/nyla-env/bin/python /home/acrossd/Desktop/nyla/nyla.py

# ── Environment ───────────────────────────────────────────
# Give the venv python priority in PATH
Environment=PATH=/home/acrossd/Desktop/nyla/nyla-env/bin:/usr/local/bin:/usr/bin:/bin
# Load the project .env file (GROQ_API_KEY etc.)
EnvironmentFile=/home/acrossd/Desktop/nyla/.env
# PulseAudio runtime socket — required for pygame/sounddevice audio
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse

# ── Restart policy ────────────────────────────────────────
# Always restart on crash with a 5-second cooldown.
# This makes the assistant self-healing — if the API blips or the
# network drops briefly, it restarts and comes back on its own.
Restart=always
RestartSec=5

# ── Logging ───────────────────────────────────────────────
# All output goes to journald — view with:
#   journalctl --user -u nyla.service -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
